# systemd-backend.service
[Unit]
Description=RentCar Backend API
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/rentcar
Environment="PATH=/var/www/rentcar/venv/bin"
ExecStart=/var/www/rentcar/venv/bin/gunicorn --workers 4 --bind 127.0.0.1:5000 --timeout 120 --worker-class sync app:app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
